# deploy/sickle.service
[Unit]
Description=Sickle — Autonomous AI Development Pipeline
After=network.target

[Service]
Type=simple
User=admin0
Group=admin0
WorkingDirectory=/home/admin0/sickle-prod
ExecStart=/home/admin0/sickle-prod/.venv/bin/python3 main.py --config config-live
EnvironmentFile=/home/admin0/sickle-prod/.env
Restart=always
RestartSec=10
StandardOutput=append:/var/log/sickle/sickle.log
StandardError=append:/var/log/sickle/sickle-error.log
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
